# ansible/roles/newrelic/templates/etc/init/newrelic.conf
description "Newrelic"
start on (started docker)
stop on runlevel [!2345]

respawn

respawn limit 100 1

# Needed to find /root/.dockercfg
env HOME="/root"
env IMAGE="indigobio/newrelic"
export HOME

pre-start script
    echo "[$(date)] Starting $JOB"
    /usr/bin/docker pull $IMAGE || true
    /usr/bin/docker rm newrelic || true
end script

exec /usr/bin/docker run \
    --name="newrelic" \
    --net="host" \
    --pid="host" \
    --env="NEW_RELIC_LICENSE_KEY={{new_relic_license_key}}" \
    --env="NEW_RELIC_SERVER_LABELS='{{new_relic_server_labels}}'" \
    --volume="/var/run/docker.sock:/var/run/docker.sock:ro" \
    --volume="/sys/fs/cgroup/:/sys/fs/cgroup:ro" \
    --volume="/dev:/dev" \
    --restart=always \
    --cap-add=SYS_PTRACE \
    --security-opt=apparmor:unconfined \
    indigobio/newrelic

post-stop script
    /usr/bin/docker stop -t 2 newrelic || true
    /usr/bin/docker rm newrelic || true
    sleep 2
end script
